# services/ceph/client.service
[Unit]
Description=Ceph Client Prerequisite
After=ceph-mds.target

[Service]
Type=forking
ExecStart=/usr/bin/ceph -s
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
