# docker/ceph/client.service
[Unit]
Description=Ceph Client Prerequisite

[Service]
Type=forking
ExecStart=/usr/bin/ceph -s
RemainAfterExit=yes
