# services/ajenti/ajenti.service
[Unit]
Description=Ajenti

[Service]
Type=forking
ExecStart=/usr/local/bin/ajenti-panel

[Install]
WantedBy=multi-user.target
